# automation/com.fomo.cursor-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fomo.cursor-agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/dom.khr/Desktop/fomofinal copy/automation/scripts/agent_daemon.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/dom.khr/Desktop/fomofinal copy/automation/logs/agent.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/dom.khr/Desktop/fomofinal copy/automation/logs/agent.error.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/dom.khr/Desktop/fomofinal copy</string>
</dict>
</plist> 